# starthere-demo/supervisord.conf
[supervisord]
nodaemon=true

[program:secret]
command=/usr/local/bin/python app.py