# foren01/supervisord.conf
[supervisord]
nodaemon=true

[program:vsftpd]
command=/bin/start_vsftpd.sh

[program:sshd]
command=/usr/sbin/sshd -D